# etc/ams-landing-restarter.path
[Path]
PathModified=${WORK_DIRECTORY}

[Install]
WantedBy=multi-user.target
